# custom_files/RPI/files_with_hotplug_sript/etc/hotplug.d/iface/99-my-action
if [ "$ACTION" = ifup ]
then
echo "$INTERFACE" >> /log_file_test
echo "$DEVICE" >> /log_file_test
fi



